# ksshaskpass/.config/autostart/ssh-add.desktop
[Desktop Entry]
Exec=ssh-add -q ~/.ssh/github ~/.ssh/badass-server ~/.ssh/personal-server
Name=ssh-add
Type=Application
